# auto-focus-cloud.service
[Unit]
Description=Auto-Focus API
After=network.target

[Service]
Type=simple
User=autofocus
WorkingDirectory=/home/autofocus/app
ExecStart=/home/autofocus/app/auto-focus-cloud
Restart=always
RestartSec=5

# Environment variables
Environment=PORT=8080
EnvironmentFile=-/home/autofocus/app/.env

# Security settings (without ProtectHome=true which was causing issues)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/autofocus/app

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto-focus-cloud

[Install]
WantedBy=multi-user.target
